# target/linux/sunxi/base-files/etc/uci-defaults/99-base-setup
#!/bin/sh

. /lib/functions.sh
board=`board_name`

uci set fstab.@global[0].anon_mount='1'
uci commit fstab

case "$board" in
	"tempe,tempe-h6v2")
		# allow lan/wan input
		uci set firewall.@zone[0].input='ACCEPT'
		uci set firewall.@zone[1].input='ACCEPT'
		uci commit firewall
		# drop network.lan
		uci delete network.lan
		uci commit network
	;;
esac
