# target/linux/sunxi/base-files/etc/init.d/reset_hw_addr
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=18

reset_hw_addr() {
	. /lib/sunxi_utils.sh
	env_reset_hw_addr
}

boot() {
	reset_hw_addr &
}
